# instguard/deploy/systemd/instguard.timer
[Unit]
Description=Run instguard supply-chain scan over the configured project periodically
Documentation=https://github.com/mtclinton/defensive-suite

[Timer]
# Every 6 hours with up to 20 min of jitter so a fleet does not stampede the
# collector or OSV.dev. A poisoned release earns a MAL- classification within
# ~3 days, so a sub-daily cadence catches a freshly-flagged pinned dependency
# well inside the cooldown window. Persistent=true runs a missed scan after wake.
OnCalendar=*-*-* 00/6:00:00
RandomizedDelaySec=20m
Persistent=true
# Also run shortly after boot.
OnBootSec=5min
Unit=instguard.service

[Install]
WantedBy=timers.target
