# authwatch/deploy/systemd/authwatch.timer
[Unit]
Description=Run authwatch trust-path integrity check daily
Documentation=https://github.com/mtclinton/defensive-suite

[Timer]
# Daily, early morning, with up to 30 min of jitter so a fleet does not stampede
# the collector. Persistent=true runs a missed check after the machine wakes.
OnCalendar=*-*-* 03:00:00
RandomizedDelaySec=30m
Persistent=true
# Also run shortly after boot — persistence often lands across a reboot.
OnBootSec=10min
Unit=authwatch.service

[Install]
WantedBy=timers.target
